# launchd/com.jarvis.voice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Edit the two /PATH/TO placeholders, copy to ~/Library/LaunchAgents/, then:
     launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.jarvis.voice.plist

     IMPORTANT (macOS quirks discovered the hard way):
     1. Invoke the VENV python directly (not bash -c 'cd ...', not the bare
        framework python) — bash cd onto an external volume hangs on TCC,
        and the framework python has no site-packages.
     2. If the project lives on an external drive, grant Full Disk Access to
        Python.app (System Settings > Privacy & Security > Full Disk Access >
        .../Python.framework/Versions/3.x/Resources/Python.app).
     3. Log paths MUST be on the internal disk or launchd fails with EX_CONFIG.
     4. The FD limits below are REQUIRED: launchd's default is 256 and the
        server idles at ~244 open files (Python libs + the Whisper model) —
        without these keys it dies with "Too many open files" within hours. -->
<plist version="1.0"><dict>
  <key>Label</key><string>com.jarvis.voice</string>
  <key>SoftResourceLimits</key><dict><key>NumberOfFiles</key><integer>8192</integer></dict>
  <key>HardResourceLimits</key><dict><key>NumberOfFiles</key><integer>8192</integer></dict>
  <key>ProgramArguments</key><array>
    <string>/PATH/TO/voice-pipeline/.venv/bin/python</string>
    <string>/PATH/TO/voice-pipeline/server.py</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/YOUR_USER/Library/Logs/jarvis-voice.log</string>
  <key>StandardErrorPath</key><string>/Users/YOUR_USER/Library/Logs/jarvis-voice.log</string>
  <key>EnvironmentVariables</key><dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict></plist>
