# launchd/com.jarvis.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Hermes dashboard, bound to loopback only (the HUD's TLS proxy exposes it).
     Edit YOUR_USER, copy to ~/Library/LaunchAgents/, then:
     launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.jarvis.dashboard.plist -->
<plist version="1.0"><dict>
  <key>Label</key><string>com.jarvis.dashboard</string>
  <key>ProgramArguments</key><array>
    <string>/Users/YOUR_USER/.hermes/hermes-agent/venv/bin/hermes</string>
    <string>dashboard</string>
    <string>--host</string><string>127.0.0.1</string>
    <string>--port</string><string>9119</string>
    <string>--no-open</string><string>--skip-build</string><string>--insecure</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/YOUR_USER/Library/Logs/jarvis-dashboard.log</string>
  <key>StandardErrorPath</key><string>/Users/YOUR_USER/Library/Logs/jarvis-dashboard.log</string>
</dict></plist>
